# vnc-app/supervisor/app.conf
[program:camera_app]
environment=DISPLAY=1

